# systemd/rclone-updatedb.timer
[Unit]
Description=rclone google drive daily 8h refresh

[Timer]
OnBootSec=5min
OnUnitActiveSec=8h

[Install]
WantedBy=timers.target
